# systemd/hereassistant-git-vault@.service
[Unit]
Description=HereAssistant Git credential vault for %i
After=local-fs.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/libexec/hereassistant-git-vault --unix-user %i --database /opt/hereassistant/bridge.sqlite3
LoadCredentialEncrypted=git-credentials.json:/etc/hereassistant/git-credentials/%i.json.cred
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadOnlyPaths=/opt/hereassistant/bridge.sqlite3 /etc/hereassistant/runners
ReadWritePaths=/run/hereassistant/git-vault
MemoryDenyWriteExecute=yes
RestrictAddressFamilies=AF_UNIX
LockPersonality=yes

[Install]
WantedBy=multi-user.target
